# docs/systemd/nianguichu.service
[Unit]
Description=NianGuiChu Memorial Next.js
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/nianguichu
Environment=NODE_ENV=production
EnvironmentFile=/var/www/nianguichu/.env
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
